# autostart/x11-layout-switch-release.desktop.in
[Desktop Entry]
Type=Application
Name=X11 Layout Switch Release
Comment=Switch keyboard layout on Alt+Shift or Ctrl+Shift release
Exec=@LISTENER_PATH@
X-GNOME-Autostart-enabled=true
Terminal=false
